# sheepit-client/pkg-plist
usr/local/bin/sheepit-client-latest.jar
usr/local/etc/sheepit-client.conf
usr/local/rc.d/sheepit_client